# files/supervisor.d/cobbler.ini
[program:cobblerd]
process_name=%(program_name)s
command=/usr/bin/cobblerd -F -l INFO
autostart=true
autorestart=true
startsecs=10
startretries=3
redirect_stderr=true
